# supervisord.conf
[supervisord]
nodaemon=true

[program:server]
command=npm run build --prefix server

[program:client]
command=npm run build --prefix client