# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:server]
command=npm run start --prefix server

[program:client]
command=npm run start --prefix client